# .config/systemd/user/webcam-recheck.service
[Unit]
Description=Re-test AMD ISP4 / MIPI webcam pipeline (Z13)
ConditionVirtualization=!container

[Service]
Type=oneshot
ExecStart=%h/.local/bin/webcam-recheck
