# .config/systemd/user/z13-battery-limit.service
[Unit]
Description=GZ302 battery charge limit (80%%) via z13ctl
# z13ctl talks to asus-wmi sysfs; needs the kernel module + udev rules up.
Requires=z13ctl.socket
After=z13ctl.socket z13ctl.service
ConditionVirtualization=!container

[Service]
Type=oneshot
# 80 is a reasonable ceiling for a laptop that spends most time docked.
# Override per-machine: systemctl --user edit z13-battery-limit.service
#   [Service]
#   Environment=Z13_CHARGE_LIMIT=100
Environment=Z13_CHARGE_LIMIT=80
ExecStart=/usr/bin/z13ctl batterylimit --set $Z13_CHARGE_LIMIT
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
