# .config/systemd/user/z13-high-performance.service
[Unit]
Description=ROG Z13 high-performance fan curve via z13ctl
# z13ctl writes the 8-point curve to asus-nb-wmi hwmon sysfs. Wait for the
# z13ctl daemon (socket-activated) so it's around for any RGB/lightbar peers
# that share the same session start.
Requires=z13ctl.socket
After=z13ctl.socket z13ctl.service
ConditionVirtualization=!container

[Service]
Type=oneshot
# 8 temp:speed pairs (Celsius:percent), monotonically increasing temps,
# non-decreasing speeds. Aggressive curve: fans hold 60% from idle-warm
# and saturate by 85 C, well before the Strix Halo APU hits tjmax (~95 C).
# Override per-machine: systemctl --user edit z13-high-performance.service
ExecStart=/usr/bin/z13ctl fancurve --set "40:60%,50:75%,60:82%,70:90%,80:98%,85:100%,90:100%,100:100%"
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
